# docker/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/var/log/supervisord.log
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log

[program:serverpaneld]
command=python3 /etc/serverpanel/daemon/bin/start_daemon.py
autostart=true
autorestart=true
stderr_logfile=/var/log/daemon.err.log
stdout_logfile=/var/log/daemon.out.log